# crontab
0 */4 * * * /usr/local/bin/python /guarding.py $(hostname)  >> /var/log/cron.log 2>&1

